# systemd/pingui-dev.service.example
[Unit]
Description=PINGUI dev session (venv, in-memory route monitor)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/path/to/PINGUI
ExecStart=/path/to/PINGUI/pingui.sh
# Розгортання: /path/to/PINGUI/pingui.sh --deploy
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
